# deploy/msgflow-listing.service
[Unit]
Description=crypto-msgflow listing announcement watcher (OKX/Binance -> Lark)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/crypto-msgflow
Environment=PYTHONUNBUFFERED=1
Environment=HTTP_PROXY=http://127.0.0.1:7890
Environment=HTTPS_PROXY=http://127.0.0.1:7890
ExecStart=/path/to/venv/bin/python3 -m watchers.listing_alert
Restart=always
RestartSec=10
StandardOutput=append:/path/to/crypto-msgflow/logs/listing.log
StandardError=append:/path/to/crypto-msgflow/logs/listing.log

[Install]
WantedBy=default.target
